# crawling_security/supervisord_scheduler.conf
[supervisord]
nodaemon=true

[program:security_scheduler]
command=/bin/bash -l -c 'cd /app; /usr/local/bundle/bin/bundle exec rake versioneye:scheduler_security_prod'
autorestart=true
